# live/config/includes.chroot/etc/systemd/system/fnos-rescue-web.service
[Unit]
Description=FNOS Rescue local web console
After=network.target

[Service]
Type=simple
User=root
Environment=PYTHONPATH=/opt/fnos-rescue/src
ExecStart=/usr/bin/python3 -m fnos_rescue serve --host 127.0.0.1 --port 8790 --static-dir /opt/fnos-rescue/web/dist
Restart=on-failure
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=/root/.local/share/fnos-rescue /run /tmp

[Install]
WantedBy=multi-user.target
